# pkg/LaunchDaemons/com.orchardpatch.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.orchardpatch.agent</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/orchardpatch/agent/run.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/var/log/orchardpatch/agent.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/orchardpatch/agent.error.log</string>

    <key>WorkingDirectory</key>
    <string>/usr/local/orchardpatch/agent</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>ORCHARDPATCH_HOME</key>
        <string>/etc/orchardpatch</string>
    </dict>

    <!-- Only accessible from localhost -->
    <key>Sockets</key>
    <dict>
        <key>Listeners</key>
        <dict>
            <key>SockFamily</key>
            <string>IPv4</string>
            <key>SockNodeName</key>
            <string>127.0.0.1</string>
            <key>SockServiceName</key>
            <string>47652</string>
        </dict>
    </dict>
</dict>
</plist>
